# scripts/hydro_grapper.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
HOME=/

# run the git.sh script and do the basic setup
* 1 * * * root /var/www/hypropo/scripts/hydro_killandstart.sh >/dev/null 2>&1
# pewpew
